# archives/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Exécution de get-keyword-from-rss.py toutes les 2h de 6h00 à 22h00 (6,8,10,12,14,16,18,20,22)
0 6-22/2 * * * root cd /app && python3 scripts/get-keyword-from-rss.py 2>&1 | tee -a /app/rapports/cron_get_keyword.log
# Vérification automatique du cron toutes les 10 minutes
*/10 * * * * root cd /app && python3 scripts/check_cron_health.py 2>&1 | tee -a /app/rapports/cron_health.log
# Exécute le scheduler chaque lundi à 6h du matin
0 6 * * 1 root cd /app && python3 scripts/scheduler_articles.py 2>&1 | tee -a /app/rapports/cron_scheduler.log
# Génère le radar thématique le dernier jour du mois à 5h00
0 5 28-31 * * root [ "$(date -d tomorrow +\%d)" = "01" ] && cd /app && python3 scripts/radar_wudd.py 2>&1 | tee -a /app/rapports/cron_radar.log
# Convertit les articles-from-rss en Markdown (nouveau format) le dernier jour du mois à 5h30
30 5 28-31 * * root [ "$(date -d tomorrow +\%d)" = "01" ] && cd /app && python3 scripts/articles_rss_to_markdown.py 2>&1 | tee -a /app/rapports/cron_rss_markdown.log
# Génère le rapport quotidien Top 10 entités des 48h à 23h00 (après la dernière collecte RSS de 22h00)
0 23 * * * root cd /app && python3 scripts/generate_48h_report.py 2>&1 | tee -a /app/rapports/cron_48h_report.log


