# archives/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Surveillance round-robin des flux RSS WUDD.opml : 1 flux par exécution, toutes les 5 minutes
# Maintient 48-heures.json à jour, puis rafraîchit la timeline, l'analyse croisée et les temps de lecture (calculs locaux, < 1 s)
*/5 * * * * root cd /app && { python3 scripts/flux_watcher.py 2>&1 | tee -a /app/rapports/cron_flux_watcher.log; python3 scripts/entity_timeline.py >> /app/rapports/cron_flux_watcher.log 2>&1; python3 scripts/cross_flux_analysis.py >> /app/rapports/cron_flux_watcher.log 2>&1; python3 scripts/enrich_reading_time.py >> /app/rapports/cron_flux_watcher.log 2>&1; }
# Surveillance des sources web sans RSS (sitemap) : toutes les 2h — max 5 articles par source par run
# Complète 48-heures.json et data/articles-from-rss/ de façon isolée, sans impacter le pipeline RSS
0 */2 * * * root cd /app && python3 scripts/web_watcher.py 2>&1 | tee -a /app/rapports/cron_web_watcher.log
# Exécution de get-keyword-from-rss.py toutes les 2h de 6h00 à 22h00 (6,8,10,12,14,16,18,20,22)
# (reconstruit 48-heures.json intégralement — filet de sécurité + enrichissement sur toute la fenêtre 7j)
0 6-22/2 * * * root cd /app && python3 scripts/get-keyword-from-rss.py 2>&1 | tee -a /app/rapports/cron_get_keyword.log
# Vérification automatique du cron toutes les 10 minutes
*/10 * * * * root cd /app && python3 scripts/check_cron_health.py 2>&1 | tee -a /app/rapports/cron_health.log
# Exécute le scheduler chaque lundi à 6h du matin
0 6 * * 1 root cd /app && python3 scripts/scheduler_articles.py 2>&1 | tee -a /app/rapports/cron_scheduler.log
# Génère le radar thématique le dernier jour du mois à 5h00
0 5 28-31 * * root [ "$(date -d tomorrow +\%d)" = "01" ] && cd /app && python3 scripts/radar_wudd.py 2>&1 | tee -a /app/rapports/cron_radar.log
# Convertit les articles-from-rss en Markdown (nouveau format) le dernier jour du mois à 5h30
30 5 28-31 * * root [ "$(date -d tomorrow +\%d)" = "01" ] && cd /app && python3 scripts/articles_rss_to_markdown.py 2>&1 | tee -a /app/rapports/cron_rss_markdown.log
# Génère le rapport quotidien Top 10 entités des 48h à 23h00 (après la dernière collecte RSS de 22h00)
0 23 * * * root cd /app && python3 scripts/generate_48h_report.py 2>&1 | tee -a /app/rapports/cron_48h_report.log
# Détecte les tendances et génère data/alertes.json chaque matín à 7h00
0 7 * * * root cd /app && python3 scripts/trend_detector.py 2>&1 | tee -a /app/rapports/cron_trends.log
# Morning Digest quotidien à 7h30 (30 min après trend_detector pour avoir alertes.json à jour)
30 7 * * * root cd /app && python3 scripts/generate_morning_digest.py --ai 2>&1 | tee -a /app/rapports/cron_morning_digest.log
# Notes de lecture quotidiennes à 8h00 (articles des 24h par tag, dans rapports/markdown/_WUDD.AI_/)
0 8 * * * root cd /app && python3 scripts/generate_reading_notes.py 2>&1 | tee -a /app/rapports/cron_reading_notes.log
# Backup des données : data/ → BACKUP_L1 → BACKUP_L2 chaque nuit à 1h00
# (avant les enrichissements NER/sentiment, après la dernière collecte de la veille)
0 1 * * * root cd /app && python3 scripts/backup_data.py 2>&1 | tee -a /app/rapports/cron_backup.log

# Enrichissement NER Round-Robin : 1 fichier/source par jour à 2h00
# Ajoute le champ "entities" aux articles qui n'en ont pas encore (tous les répertoires data/)
0 2 * * * root cd /app && python3 scripts/enrich_entities.py 2>&1 | tee -a /app/rapports/cron_enrich_entities.log
# Enrichissement images : fetch og:image pour les articles sans images à 2h30
# Parcourt data/articles/ et data/articles-from-rss/ — aucun appel IA, fetch HTTP uniquement
30 2 * * * root cd /app && python3 scripts/enrich_images.py 2>&1 | tee -a /app/rapports/cron_enrich_images.log
# Enrichissement sentiment Round-Robin : 1 fichier/source par jour à 3h00
# Tourne sur tous les fichiers data/articles/ et data/articles-from-rss/ de façon cyclique
0 3 * * * root cd /app && python3 scripts/enrich_sentiment.py 2>&1 | tee -a /app/rapports/cron_sentiment.log
# Réparation des résumés en erreur : chaque dimanche à 4h00 (maintenance hebdomadaire)
# Régénère les articles dont le Résumé contient un message d'erreur API EurIA
0 4 * * 0 root cd /app && python3 scripts/repair_failed_summaries.py 2>&1 | tee -a /app/rapports/cron_repair.log
# Briefing exécutif hebdomadaire : chaque lundi à 6h30 (après trend_detector du lundi matin)
# Synthèse IA : top entités, articles scorés, tendances — rapports/markdown/_BRIEFING_/
30 6 * * 1 root cd /app && python3 scripts/generate_briefing.py --period weekly 2>&1 | tee -a /app/rapports/cron_briefing.log
# Rapports mensuels par mot-clé : dernier jour du mois à 6h00
# Génère un rapport Markdown par fichier data/articles-from-rss/ pour la période courante
0 6 28-31 * * root [ "$(date -d tomorrow +\%d)" = "01" ] && cd /app && python3 scripts/generate_keyword_reports.py 2>&1 | tee -a /app/rapports/cron_keyword_reports.log

