# crontab
# Vérification automatique du cron toutes les 10 minutes
*/10 * * * * root cd /app && python3 scripts/check_cron_health.py 2>&1 | tee -a /app/rapports/cron_health.log


SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Exécute le scheduler chaque lundi à 6h du matin
0 6 * * 1 root cd /app && python3 scripts/scheduler_articles.py 2>&1 | tee -a /app/rapports/cron_scheduler.log

# Tâche de test : écrit un fichier toutes les minutes
* * * * * root echo "cron ok $(date)" | tee -a /app/cron_test.log


